# saltstack/salt/ssh_config/files/tunnels.ini
[program:redshift_tunnel]
command=autossh -N redshift-dev
user=vagrant
autostart=true
autorestart=true
stopsignal=KILL